# jackett/root/etc/services.d/jackett/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Runs Jackett
# ==============================================================================
bashio::log.info "Starting Jackett..."

if bashio::config.true 'openvpn_enabled'; then
    exec /usr/sbin/openvpn --config /etc/openvpn/config.ovpn --script-security 2 --up /etc/openvpn/up-jackett.sh
else
    #exec s6-setuidgid jackett mono /app/Jackett/JackettConsole.exe
    exec s6-setuidgid abc /app/Jackett/jackett --NoUpdates
fi
