# jackett/root/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================

# Wait for jackett to become available
bashio::net.wait_for 9117

#bashio::log.info "Starting NGinx..."

#exec nginx
